# local_deploy/2/deploy_production/scripts/astroshield-backend 2.service
[Unit]
Description=AstroShield Backend Service
After=network.target

[Service]
User=stardrive
Group=stardrive
WorkingDirectory=/home/stardrive/astroshield/backend
ExecStart=/home/stardrive/astroshield/backend/venv/bin/python3 minimal_server.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
Environment="PYTHONPATH=/home/stardrive/astroshield"
Environment="PRODUCTION=true"
Environment="HOST=0.0.0.0"
Environment="PORT=3001"

[Install]
WantedBy=multi-user.target